# supervisord.conf
[supervisord]
nodaemon=true

[program:bbb]
command=/usr/bin/start.sh
